# arch/.config/systemd/user/checkaur.service
[Unit]
Description=Run a command to poll for AUR updates

[Service]
Type=oneshot
ExecStart=%h/dotfiles/arch/checkaur
