# autostart/.config/autostart/redshift.desktop
[Desktop Entry]
Name=Redshift
Comment=
Exec=/usr/bin/redshift -l geoclue2
Terminal=false
Type=Application

